# supervisor/junior.conf
[program:junior-web]
user = ubuntu
command = /opt/env/junior/bin/gunicorn junior.wsgi:application
environment=PYTHONPATH="/opt/junior"