# supervisor/junior.conf
[program:junior-web]
command = /opt/env/junior/bin/gunicorn junior.wsgi:application